# src/config/supervisor/conf.d/npcd.conf
[program:npcd]
command=/usr/local/pnp4nagios/bin/npcd -d -f /usr/local/pnp4nagios/etc/npcd.cfg
autostart=true
user=root
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor_npcd.log
