# src/config/supervisor/conf.d/thruk.conf
[program:thrukch]
command=/bin/chown -R www-data:www-data /etc/thruk /etc/skconf
autostart=true
user=root
log_stdout=true
log_stderr=true
logfile=/var/log/supervisor_thruk.log
